# supervisor/kafka-manager.conf
[program:kafka]
command = /usr/bin/start-kafka-manager.sh
autostart = true
autorestart = true